# .config/systemd/user/aide-check.service
[Unit]
Description=Check integrity with AIDE
After=network.target

[Service]
Type=oneshot
ExecStart=%h/.local/sbin/aide
