# .config/systemd/user/aide-check.timer
[Unit]
Description=Plan check AIDE

[Timer]
OnBootSec=10m
OnUnitActiveSec=1d
Persistent=true

[Install]
WantedBy=timers.target
